# deploy/receipt-printer.service.template
[Unit]
Description=Receipt Printer Service
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%WORKDIR%
EnvironmentFile=%WORKDIR%/.env
ExecStart=%PYTHON% %WORKDIR%/app.py --server
Restart=always
RestartSec=5
LimitNOFILE=4096

# Ensure clean restarts
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
