# deploy/receipt-printer-healthcheck.service.template
[Unit]
Description=Receipt Printer Health Check
After=receipt-printer.service

[Service]
Type=oneshot
ExecStart=%WORKDIR%/scripts/healthcheck
