# deploy/receipt-printer-healthcheck.timer.template
[Unit]
Description=Run Receipt Printer Health Check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min

[Install]
WantedBy=timers.target
